# relay/systemd/punkto-relay.service
[Unit]
Description=punkto-relay v0.1 — Flow TV relay node
Documentation=https://github.com/your-org/punkto/blob/main/punkto.relay.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=punkto
Group=punkto
WorkingDirectory=/opt/punkto-relay
ExecStart=/usr/bin/python3 /opt/punkto-relay/relay.py
EnvironmentFile=-/etc/punkto-relay/relay.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/punkto-relay
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
